# contrib/systemd/user/scrbblr.service
[Unit]
Description=MPRIS Scrobbler watcher
Documentation=https://github.com/arturmeski/scrbblr
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/scrbblr watch --player com.blitzfc.qbz
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
